# .archived-code/runtime-24-7/blackbox5.service
[Unit]
Description=BLACKBOX5 24/7 Autonomous AI Runtime
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=shaansisodia
Group=staff
WorkingDirectory=/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL/blackbox5
Environment="PYTHONPATH=/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL/blackbox5"
Environment="ANTHROPIC_API_KEY=%h"
ExecStart=/usr/local/bin/python3 runtime/launcher.py
Restart=always
RestartSec=10
StandardOutput=append:/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL/blackbox5/blackbox5_runtime.log
StandardError=append:/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL/blackbox5/blackbox5_error.log

[Install]
WantedBy=multi-user.target
