# 6-roadmap/research/external/YouTube/AI-Improvement-Research/deploy/rss-scraper.service
[Unit]
Description=YouTube RSS Scraper
After=network.target

[Service]
Type=oneshot
User=transcript
WorkingDirectory=/opt/transcripts/repo
Environment="PYTHONPATH=/opt/transcripts/repo/scripts"
ExecStart=/opt/transcripts/repo/venv/bin/python scripts/queue/populate.py

[Install]
WantedBy=multi-user.target
