# 6-roadmap/research/external/YouTube/AI-Improvement-Research/deploy/rss-scraper.timer
[Unit]
Description=Run RSS Scraper every 4 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
