# agents/moltbot-autonomous/bb5-reporter.service
[Unit]
Description=BB5 Agent Activity Reporter
After=network.target moltbot-autonomous.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/blackbox5
Environment="HOME=/root"
Environment="TELEGRAM_BOT_TOKEN=8581639813:AAFA13wDTKEX2x6J-lVfpq9QHnsGRnB1EZo"
Environment="TELEGRAM_CHAT_ID=7643203581"
Environment="PYTHONPATH=/opt/blackbox5/agents/moltbot-autonomous"
ExecStart=/usr/bin/python3 /opt/blackbox5/agents/moltbot-autonomous/bb5-reporter.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
